# backend/crontab
SHELL=/bin/bash
PATH=/usr/bin:/bin:/usr/local/bin
*/1 * * * * /bin/bash /app/scripts/daily_scores.sh >> /var/log/cron.log 2>&1
